# deploy/droplet/comfytv-stream.service
[Unit]
Description=ComfyTV stream proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/comfytv-stream/server.js
EnvironmentFile=/etc/comfytv-stream/.env
Restart=always
RestartSec=3
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
